# deploy/hxz-fortune.service
[Unit]
Description=HXZ Fortune Gunicorn
Documentation=https://github.com/kiosk/hxz-fortune
After=network.target

[Service]
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=/opt/hxz-fortune-telling/backend
Environment="PATH=/opt/hxz-fortune-telling/backend/venv/bin"

ExecStart=/opt/hxz-fortune-telling/backend/venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:8002 \
    --timeout 60 \
    --access-logfile /var/log/hxz-fortune/access.log \
    --error-logfile /var/log/hxz-fortune/error.log \
    app:app

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
